# docker/root/etc/s6-overlay/s6-rc.d/razord/run
#!/bin/sh
# razord longrun — persistent Razor daemon (vendored razorfy). Keeps a Razor2
# agent warm so the shim avoids the per-message perl-startup cost of razor-check.
# Runs as the unprivileged drp user (owns /var/lib/razor). Loopback only.
exec 2>&1
export RAZORFY_BINDADDRESS=127.0.0.1
export RAZORFY_BINDPORT=11342
export RAZORFY_RAZORHOME=/var/lib/razor
export HOME=/var/lib/razor
exec s6-setuidgid drp perl /opt/drp/razord/razorfy.pl
