# root/etc/s6-overlay/s6-rc.d/shim/run
#!/bin/sh
# spamcheck_shim longrun — HTTP wrapper around dccproc / razor-check / pyzor that
# rspamd's dcc_razor_pyzor.lua plugin queries. Foreground; s6 supervises.
exec 2>&1
export RAZORHOME="${RAZORHOME:-/var/lib/razor}"
export PYZOR_HOME="${PYZOR_HOME:-/var/lib/pyzor}"
exec python3 /opt/drp/shim/spamcheck_shim.py
