# root/etc/s6-overlay/s6-rc.d/dccifd/run
#!/bin/sh
# dccifd longrun — the DCC interface daemon (from the `dcc` package). Runs in the
# foreground as the dcc user; s6 respawns it on exit. `-b` keeps it attached to
# the tty (no self-backgrounding) so s6 can supervise it. If the `dcc` package is
# absent, idle instead of crash-looping.
exec 2>&1
if [ ! -x /var/dcc/libexec/dccifd ]; then
    echo "[dccifd] /var/dcc/libexec/dccifd absent — DCC disabled; idling" >&2
    exec sleep infinity
fi
exec s6-setuidgid dcc /var/dcc/libexec/dccifd -b -h /var/dcc
